# turku-storage-update-config.service
[Unit]
Description=turku-storage-update-config

[Service]
Type=oneshot
ExecStart=/usr/bin/env turku-storage-update-config
