# turku-storage.cron
*/5 * * * * root turku-storage-update-config --wait=300
